# mr_movies.service
[Unit]
Description=Mr. Will Movies Discord Bot

[Service]
ExecStart=/bin/node /opt/mr_movies/app.js

[Install]
WantedBy=multi-user.target